# files/crontabs/root
#
# m  h  dom  mon  dow    usercommand
*/15    *  *    *    *   /usr/bin/dsl-stats.sh

